# .scripts/backup/crontab
# At 13:00 on Monday.
0 13 * * 1 $HOME/.dotfiles/.scripts/backup/mailOnErr.sh $HOME/.dotfiles/.scripts/backup/lastpass.sh
